# crontab/20190729_09261564385204.cron
*/1 * * * * Rscript /Users/federicopfahler/Desktop/scripts/r_scripts_test/test.R > /Users/federicopfahler/Desktop/scripts/logs/test.R_temp.log 2>&1 || cat /Users/federicopfahler/Desktop/scripts/logs/test.R_temp.log && cat /Users/federicopfahler/Desktop/scripts/logs/test.R_temp.log >> /Users/federicopfahler/Desktop/scripts/logs/test.R20190729_09261564385204_persistent.log && rm /Users/federicopfahler/Desktop/scripts/logs/test.R_temp.log
